# contrib/systemd-credentials-vault.service
[Unit]
Description=Systemd Credentials Vault Server
Documentation=https://github.com/midoBB/systemd-credentials-vault
After=network-online.target
Wants=network-online.target
Requires=systemd-credentials-vault.socket

[Service]
Type=exec
ExecStart=/usr/local/bin/systemd-credentials-vault -config /usr/local/etc/systemd-credentials-vault/config.yml
Restart=on-failure
RestartSec=5
TimeoutStartSec=30
TimeoutStopSec=30

# Load the Vault secret ID credential
LoadCredentialEncrypted=vault-secret-id:/etc/vault-server/secret-id.cred

[Install]
WantedBy=multi-user.target
Also=systemd-credentials-vault.socket
